# docker/systemd/gaeld-worker.service.example
[Unit]
Description=Gaeld Queue Worker
After=network.target redis.service postgresql.service

[Service]
User=deploy
Group=www-data
WorkingDirectory=/var/www/gaeld/current
ExecStart=/usr/bin/php /var/www/gaeld/current/artisan queue:work redis \
    --sleep=3 \
    --tries=3 \
    --max-time=3600 \
    --queue=default
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=3600
StandardOutput=append:/var/www/gaeld/shared/storage/logs/worker.log
StandardError=append:/var/www/gaeld/shared/storage/logs/worker.log

[Install]
WantedBy=multi-user.target
